# sysadmin/prod/roles/python/files/wynno_python.conf
start on runlevel [2345]
stop on runlevel [!2345]

script
    python /home/deploy/dist/lib/python/server.py
end script
