# sysadmin/prod/roles/python/files/wynno_python.conf
description "wynno Python machine learning"
author  "engineering@wynno.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec python /home/deploy/dist/lib/python/server.py
